# docker/supervisord-api.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

; ── nginx（前端静态 + API 反代）─────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── API Server ───────────────────────────────────────────────
[program:fanapi-server]
command=/app/fanapi-server
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=3
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
